# container-files/etc/services.d/exim/run
#!/bin/execlineb -P

/usr/sbin/exim -bdf
